# liga_edumind_backend.service
[Unit]
Description=Liga EDUmind Backend API (FastAPI)
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=nuevoadmin
Group=nuevoadmin
WorkingDirectory=/var/www/liga_edumind/backend
Environment="PATH=/var/www/liga_edumind/backend/venv/bin"
ExecStart=/var/www/liga_edumind/backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8001 --workers 4
Restart=always
RestartSec=10

#  Hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
